# jarvis.service.example
[Unit]
Description=Jarvis Assistant
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/jarvis
ExecStart=/path/to/uv run python -m jarvis.main
Restart=always
RestartSec=3
EnvironmentFile=/path/to/jarvis/.env

[Install]
WantedBy=multi-user.target
